# newsboat/.config/systemd/user/newsup.service
[Unit]
Description=Newsboat RSS Feed Update Daemon

[Service]
ExecStart=%h/.local/bin/newsup
