# newsboat/.config/systemd/user/newsup.service
[Unit]
Description=Sync Rss feeds

[Service]
Type=oneshot
ExecStart=newsboat -x reload
