# newsboat/.config/systemd/user/newsup.timer
[Unit]
Description=Update rss feeds timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=10m
Unit=newsup.service

[Install]
WantedBy=timers.target
