# borg/.config/systemd/user/borg.service
[Unit]
Description=Borg Backup

[Service]
ExecStart=/home/alrayyes/bin/backup.sh

[Install]
WantedBy=default.target
