# slstatus/.config/systemd/user/update-pacman-cache.service
[Unit]
Description=Update Pacman Cache

[Service]
ExecStart=yay -Sy

[Install]
WantedBy=default.target
