# slstatus/.config/systemd/user/update-pacman-cache.timer
[Unit]
Description=Run update-pacman-cache every 15m

[Timer]
OnBootSec=300m
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target
