# restic/.config/systemd/user/restic.timer
[Unit]
Description=Run restic daily

[Timer]
OnBootSec=300m
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
