# restic/.config/systemd/user/restic.service
[Unit]
Description=Restic Backup

[Service]
ExecStart=/home/alrayyes/.scripts/backup.sh

[Install]
WantedBy=default.target
