# restic/.config/systemd/user/restic.service
[Unit]
Description=Restic Backup

[Service]
ExecStart=/home/alrayyes/.local/bin/backup

[Install]
WantedBy=default.target
